# taxi/crontab
*/5 * * * * /usr/bin/taxi.sh to_office >/proc/1/fd/1 2>/proc/1/fd/2
*/5 * * * * /usr/bin/taxi.sh from_office >/proc/1/fd/1 2>/proc/1/fd/2

